# systemd/user/elephant-hyprland.service

[Unit]
Description=Elephant
After=xdg-desktop-autostart.target

[Service]
ExecStart=elephant
Restart=on-failure
